# kde.conf
[program:kde]
directory=/home/user
command=/usr/bin/startkde
user=user
autostart=true
autorestart=true
stopsignal=QUIT
environment=DISPLAY=":0",HOME="/home/user",SHELL="/bin/bash"
